# shadow_assists/features/IDE_015_cpu_extreme_util/SUB_165_cgroup_isolcpus_doc/cgroup_yaml/vllm_trident.slice
[Unit]
Description=vLLM trident backend (GPU 4-7, NUMA 1)
PartOf=spec_decode.slice
Before=vllm_trident.service

[Slice]
AllowedCPUs=62-79
AllowedMemoryNodes=1
CPUWeight=200
MemoryHigh=512G
