# shadow_assists/features/IDE_015_cpu_extreme_util/SUB_165_cgroup_isolcpus_doc/cgroup_yaml/vllm_vanilla.slice
[Unit]
Description=vLLM vanilla backend (GPU 0-3, NUMA 0)
PartOf=spec_decode.slice
Before=vllm_vanilla.service

[Slice]
AllowedCPUs=6-30
AllowedMemoryNodes=0
CPUWeight=200
MemoryHigh=512G
